# ops/systemd/fjcloud-api.service
# /etc/fjcloud/env is NOT tracked in git. The deploy script runs
# generate_ssm_env.sh to populate it from SSM parameters under /fjcloud/<env>/*
# before each service restart. Mapping: ops/scripts/lib/generate_ssm_env.sh.
# File permissions: 0600 fjcloud:fjcloud.
# Production startup contract: when ENVIRONMENT is prod/production,
# /etc/fjcloud/env must include at least one non-blank webhook:
# SLACK_WEBHOOK_URL or DISCORD_WEBHOOK_URL.

[Unit]
Description=fjcloud API Server
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=fjcloud
Group=fjcloud

EnvironmentFile=-/etc/fjcloud/env
ExecStart=/usr/local/bin/fjcloud-api

Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fjcloud-api

[Install]
WantedBy=multi-user.target
